# systemd/tgbot_recipes.service
[Unit]
Description=Bot Recipes
After=multi-user.target

[Service]
User=root
Group=admin
Type=simple
WorkingDirectory=/root/tgbot_recipes
ExecStart=/root/tgbot_recipes/venv/bin/python /root/tgbot_recipes/bot.py
Restart=always

[Install]
WantedBy=multi-user.target
